# spaceship/spaceship.service
[Unit]
Description=Spaceship XPlanet Server
After=network.target

[Service]
Type=simple
User=atacama
Group=atacama
WorkingDirectory=/home/atacama/atacama
Environment=PYTHONPATH=/home/atacama/atacama
ExecStart=/usr/bin/python3 /home/atacama/atacama/launch.py --mode spaceship
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
